# systemd/kitty-session-save.service
[Unit]
Description=Save Kitty terminal session

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/adj/.config/kitty/save_session.py
